# conf/oxe_spot.service
[Unit]
Description=Oxe Spot Daemon

[Service]
WorkingDirectory=/home/pi/oxe_spot
ExecStart=/home/pi/oxe_spot/scripts/oxe_spotd.sh

[Install]
WantedBy=default.target
